# service/docker-rbd-volume-driver.service.d/rbd.conf
[Service]
Environment="PATH=/opt/bin:/usr/sbin:/usr/bin"
